# supervisor.conf
[supervisord]
nodaemon=true

[program:redis]
priority=10
directory=/tmp
command=redis-server
user=root
autostart=true
autorestart=true

[program:hubot]
priority=20
command=/opt/hubot/bin/hubot --adapter flowdock --name mbot
user=root
autostart=true
autorestart=true